# scripts/arpit-brain.service
[Unit]
Description=Arpit Brain Trainer v2.0 — God Mode Autonomous Agent
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/arpit-brain-trainer
ExecStart=/usr/bin/node src/brain.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
EnvironmentFile=/home/ec2-user/arpit-brain-trainer/.env

# Logging
StandardOutput=append:/home/ec2-user/arpit-brain-trainer/logs/brain.log
StandardError=append:/home/ec2-user/arpit-brain-trainer/logs/brain-error.log

# Resource limits
MemoryMax=512M
CPUQuota=80%

[Install]
WantedBy=multi-user.target
